# build/linux/fpm.before-remove.sh
#!/usr/bin/env bash

systemctl stop skyhwd.service
systemctl disable skyhwd.service
