# build/linux/fpm.after-install.sh
#!/bin/bash

systemctl enable skyhwd.service
systemctl start skyhwd.service
